# src/templates/wifi-fallback.service
[Unit]
Description=WiFi Fallback Service

[Service]
ExecStart=/usr/local/bin/wifi-fallback.sh

