# src/templates/wifi-fallback.timer
[Unit]
Description=WiFi Fallback Timer

[Timer]
OnBootSec=30s
OnUnitActiveSec=15s

[Install]
WantedBy=timers.target

